# roles/vdirsyncer_container/files/container-vdirsyncer.timer
[Unit]
Description=Execute vdirsyncer container daily

[Timer]
OnCalendar=*-*-* 5:00:00

[Install]
WantedBy=timers.target
