# roles/autorestic/files/autorestic.timer
[Unit]
Description=Daily autorestic backup

[Timer]
OnCalendar=*-*-* 6:00:00

[Install]
WantedBy=timers.target
