# roles/transmission/templates/gluetun.service.j2
# gluetun.service
# autogenerated by Podman 4.7.0
[Unit]
Description=Gluetun
Documentation=man:podman-generate-systemd(1)
Wants=network-online.target
After=network-online.target
RequiresMountsFor=%t/containers
BindsTo=transmission-gluetun.service
After=transmission-gluetun.service

[Service]
Environment=PODMAN_SYSTEMD_UNIT=%n
Restart=on-failure
TimeoutStopSec=70
ExecStart=/usr/bin/podman container run \
	--cidfile=%t/%n.ctr-id \
	--cgroups=no-conmon \
	--rm \
	--pod-id-file %t/transmission-gluetun.pod-id \
	--sdnotify=conmon \
	-d \
	--replace \
	--name gluetun \
	--cap-add NET_ADMIN \
	--sysctl net.ipv6.conf.all.disable_ipv6=0 \
	--device /dev/net/tun:/dev/net/tun \
	--volume {{ containers_storage_path }}/data/gluetun:/gluetun:Z \
	--tz={{ timezone }} \
	--env-file {{ containers_storage_path }}/env/gluetun/gluetun.env \
	--env TZ={{ timezone }} \
	--env VPN_TYPE=wireguard \
	--env VPN_INTERFACE=wg0 \
	--env VPN_SERVICE_PROVIDER={{ gluetun.provider.name }} \
	--env VPN_ENDPOINT_IP={{ gluetun.wireguard.endpoint.ip }} \
	--env VPN_ENDPOINT_PORT={{ gluetun.wireguard.endpoint.port }} \
	--env WIREGUARD_PUBLIC_KEY={{ gluetun.wireguard.endpoint.public_key }} \
	--env WIREGUARD_ADDRESSES={{ gluetun.wireguard.addresses }} \
	--env FIREWALL_VPN_INPUT_PORTS={{ gluetun.firewall.vpn_input_ports }} \
	--env DOT=off \
	--env DNS_ADDRESS={{ gluetun.dns.addresses }} {{ containers_images.gluetun }}:{{ containers_images_tags.gluetun }}
ExecStop=/usr/bin/podman stop \
	--ignore -t 10 \
	--cidfile=%t/%n.ctr-id
ExecStopPost=/usr/bin/podman rm \
	-f \
	--ignore -t 10 \
	--cidfile=%t/%n.ctr-id
Type=notify
NotifyAccess=all

[Install]
WantedBy=default.target
