# roles/transmission/templates/transmission.service.j2
# transmission.service
# autogenerated by Podman 4.7.0
[Unit]
Description=Transmission
Documentation=man:podman-generate-systemd(1)
Wants=network-online.target
After=network-online.target
RequiresMountsFor=%t/containers
BindsTo=transmission-gluetun.service
After=transmission-gluetun.service

[Service]
Environment=PODMAN_SYSTEMD_UNIT=%n
Restart=on-failure
TimeoutStopSec=70
ExecStart=/usr/bin/podman container run \
	--cidfile=%t/%n.ctr-id \
	--cgroups=no-conmon \
	--rm \
	--pod-id-file %t/transmission-gluetun.pod-id \
	--sdnotify=conmon \
	-d \
	--replace \
	--name transmission \
	--expose 9091 \
	--volume {{ containers_storage_path }}/data/transmission:/config:Z \
	--volume {{ nfs_mounts_path }}/downloads/bittorrent/downloads:/downloads:slave \
	--volume {{ nfs_mounts_path }}/downloads/bittorrent/watch:/watch:slave \
	--secret transmission_user \
	--secret transmission_password \
	--tz={{ timezone }} \
	--env TZ={{ timezone }} \
	--env FILE__USER=/run/secrets/transmission_user \
	--env FILE__PASS=/run/secrets/transmission_password \
	--label traefik.enable=true \
	--label traefik.docker.network=common-internal \
	--label traefik.http.routers.transmission.rule=Host(`transmission.{{ domain_name }}`) \
	--label traefik.http.routers.transmission.entrypoints=websecure \
	--label "traefik.http.routers.transmission.middlewares=set-security-headers@file, restrict-external-access@file" \
	--label traefik.http.services.transmission.loadbalancer.server.port=9091 {{ containers_images.transmission }}:{{ containers_images_tags.transmission }}
ExecStop=/usr/bin/podman stop \
	--ignore -t 10 \
	--cidfile=%t/%n.ctr-id
ExecStopPost=/usr/bin/podman rm \
	-f \
	--ignore -t 10 \
	--cidfile=%t/%n.ctr-id
Type=notify
NotifyAccess=all

[Install]
WantedBy=default.target
