# roles/transmission/templates/transmission-gluetun.service.j2
# transmission-gluetun.service
# autogenerated by Podman 4.7.0
[Unit]
Description=Tranmission-Gluetun Pod
Documentation=man:podman-generate-systemd(1)
Wants=network-online.target
After=network-online.target
RequiresMountsFor=/run/containers/storage
Wants=gluetun.service transmission.service
Before=gluetun.service transmission.service

[Service]
Environment=PODMAN_SYSTEMD_UNIT=%n
Restart=on-failure
TimeoutStopSec=70
ExecStartPre=/usr/bin/podman pod create \
	--infra-conmon-pidfile %t/transmission-gluetun.pid \
	--pod-id-file %t/transmission-gluetun.pod-id \
	--exit-policy=stop \
	--name transmission-gluetun \
	--network gluetun-external \
	--network common-internal \
	--replace
ExecStart=/usr/bin/podman pod start \
	--pod-id-file %t/transmission-gluetun.pod-id
ExecStop=/usr/bin/podman pod stop \
	--ignore \
	--pod-id-file %t/transmission-gluetun.pod-id  \
	-t 10
ExecStopPost=/usr/bin/podman pod rm \
	--ignore \
	-f \
	--pod-id-file %t/transmission-gluetun.pod-id
PIDFile=%t/transmission-gluetun.pid
Type=forking

[Install]
WantedBy=default.target
