# roles/fedora_coreos_setup/files/bootupd-auto-update.service
[Unit]
Description=Bootupd automatic update

[Service]
ExecStart=/usr/bin/bootupctl update
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
