# roles/fedora_coreos_setup/files/rpm-ostree-refresh-md.service
[Unit]
Description=rpm-ostree repository metadata refresh

[Service]
ExecStart=/usr/bin/rpm-ostree refresh-md
Type=oneshot

[Install]
WantedBy=multi-user.target
