# roles/fedora_coreos_setup/templates/systemd.automount.j2
[Unit]
Description={{ storage_server.name }} {{ item }} nfs share

[Automount]
Where={{ nfs_mounts_path }}/{{ item }}

[Install]
WantedBy=multi-user.target
