# roles/healthchecks/files/healthchecks-postgres-dump.timer
[Unit]
Description=Daily Healthchecks PostgreSQL dump

[Timer]
OnCalendar=*-*-* 4:00:00

[Install]
WantedBy=timers.target
