# roles/healthchecks_containers/files/container-exec-healthchecks-postgres-dump.timer
[Unit]
Description=Execute Healthchecks Postgres dump daily

[Timer]
OnCalendar=*-*-* 4:00:00

[Install]
WantedBy=timers.target
