# roles/fedora_iot_setup/templates/smart-test@.service.j2
[Unit]
Description=SMART long test of %f
OnSuccess=ping-healthchecks-slug@{{ ansible_hostname }}-smart-test-%i:success.service
OnFailure=ping-healthchecks-slug@{{ ansible_hostname }}-smart-test-%i:failure.service

[Service]
ExecStart=/usr/sbin/smartctl -t long %f
Type=simple
