# roles/nextcloud/files/nextcloud-cron.timer
[Unit]
Description=Execute Nextcloud cron.php every 5 minutes

[Timer]
OnCalendar=*-*-* *:0/5:00

[Install]
WantedBy=timers.target
