# roles/nextcloud/files/nextcloud-mariadb-dump.timer
[Unit]
Description=Daily Nextcloud MariaDB dump

[Timer]
OnCalendar=*-*-* 4:00:00

[Install]
WantedBy=timers.target
