# roles/nextcloud_containers/files/container-exec-nextcloud-mariadb-dump.timer
[Unit]
Description=Execute Nextcloud MariaDB dump daily

[Timer]
OnCalendar=*-*-* 4:00:00

[Install]
WantedBy=timers.target
