# roles/podman_setup/files/podman-auto-prune.timer
[Unit]
Description=Daily pruning of dangling container images

[Timer]
OnCalendar=daily
RandomizedDelaySec=900
Persistent=true

[Install]
WantedBy=timers.target
