# roles/podman_setup/files/podman-auto-prune.timer
[Unit]
Description=Podman auto-prune timer

[Timer]
OnCalendar=daily
RandomizedDelaySec=900
Persistent=true

[Install]
WantedBy=timers.target
