# roles/immich/files/immich-postgres-dump.timer
[Unit]
Description=Daily Immich PostgreSQL dump

[Timer]
OnCalendar=*-*-* 4:00:00

[Install]
WantedBy=timers.target
