# roles/rsync_servers/files/rsync-servers.timer
[Unit]
Description=Daily backups of servers data

[Timer]
OnCalendar=*-*-* 5:00:00

[Install]
WantedBy=timers.target
