# roles/fedora_server_setup/templates/btrfs-scrub@.service.j2
[Unit]
Description=Btrfs scrub on %f
ConditionPathIsMountPoint=%f
RequiresMountsFor=%f
Wants=ping-healthchecks-io-slug@{{ ansible_hostname }}-btrfs-scrub-%i:start.service
OnSuccess=ping-healthchecks-io-slug@{{ ansible_hostname }}-btrfs-scrub-%i:success.service
OnFailure=ping-healthchecks-io-slug@{{ ansible_hostname }}-btrfs-scrub-%i:failure.service

[Service]
Nice=19
IOSchedulingClass=idle
KillSignal=SIGINT
ExecStart=/usr/sbin/btrfs scrub start -B %f
